# scripts/systemd/booking-platform-api.service
[Unit]
Description=Booking Platform API
After=network.target

[Service]
Type=simple
User=booking
Group=booking
WorkingDirectory=/opt/booking-platform/backend
EnvironmentFile=/opt/booking-platform/backend/.env
ExecStart=/usr/bin/node /opt/booking-platform/backend/src/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
CapabilityBoundingSet=
PrivateTmp=true

[Install]
WantedBy=multi-user.target
